# docs/supervisor-queue.conf
[program:auflow-queue]
command=php /var/www/auflow/artisan queue:work redis --queue=default --sleep=3 --tries=3 --max-time=3600
directory=/var/www/auflow
user=www-data
numprocs=2
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/auflow-queue.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5

[program:auflow-notifications]
command=php /var/www/auflow/artisan queue:work redis --queue=notifications --sleep=3 --tries=3 --max-time=3600
directory=/var/www/auflow
user=www-data
numprocs=1
autostart=true
autorestart=true
startsecs=1
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/auflow-notifications.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
